# prom-monitor/roles/node_exporter/templates/node_exporter.service
[Unit]
Description=Node Exporter

[Service]
ExecStart={{ node_exporter_path }}/node_exporter

[Install]
WantedBy=multi-user.target
